# deploy/robertapp-backend.service
[Unit]
Description=RobertApp Backend (FastAPI + uvicorn)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Ahorros/backend
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
